# eww/.config/autostart/eww-clock.desktop
[Desktop Entry]
Type=Application
Name=eww clock
Exec=/home/amal/.config/eww/launch.sh
X-GNOME-Autostart-enabled=true
